# templates/supervisor.conf
# Gunicorn 的配置文件 /etc/supervisor/conf.d/gunicorn.conf
[program:gunicorn-<web-name>]
command = /home/.pyenvs/<web-name>/bin/gunicorn -c server-config/gunicorn-config.py <web-name>.wsgi:application
directory = /home/<web-name>
user=root
startsecs=0
stopwaitsecs=0
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/gunicorn.access.log
stderr_logfile = /var/log/supervisor/gunicorn.error.log